# motor/radar-mercado.service
[Unit]
Description=Radar Fundiario - aquecimento noturno do cache de mercado (AQUEC-01)
After=network.target postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/opt/radar/api
EnvironmentFile=/opt/radar/api/.env
ExecStart=/usr/bin/node /opt/radar/api/mercado-aquecer.js
User=www-data
NoNewPrivileges=true
PrivateTmp=true
TimeoutStartSec=2h
# Cota de CPU (alerta Hostinger 21/07): lote com prioridade baixa e teto de 1,5 nucleo.
Nice=15
CPUQuota=150%
IOWeight=10
