# motor/radar-backup.timer
[Unit]
Description=Backup diario do banco radar

[Timer]
OnCalendar=*-*-* 06:20:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
