# motor/radar-intelligence.timer
[Unit]
Description=Corretor Inteligente - investigação inteligente diária pós-varredura

[Timer]
OnBootSec=2m
OnUnitInactiveSec=10m
RandomizedDelaySec=1m
AccuracySec=30s

[Install]
WantedBy=timers.target
