# motor/radar-intelligence.service
[Unit]
Description=Corretor Inteligente - orquestração Kimi K3 de pesquisas e oportunidades
After=network-online.target postgresql.service docker.service radar-api.service

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/opt/radar/api
EnvironmentFile=/opt/radar/api/.env
ExecStart=/usr/bin/node /opt/radar/api/intelligence-orchestrator.js
NoNewPrivileges=true
PrivateTmp=true
TimeoutStartSec=3h
# Cota de CPU (alerta Hostinger 21/07): a fila roda em loop o dia todo (timer 10 min);
# quando um job dispara (plano + lotes + sintese K3), fica no ultimo nucleo disponivel.
Nice=19
CPUQuota=150%
IOWeight=1
