# motor/radar-agent-review.service
[Unit]
Description=Corretor Inteligente - revisão econômica do agente
After=postgresql.service radar-api.service

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/opt/radar/api
EnvironmentFile=/opt/radar/api/.env
ExecStart=/usr/bin/node /opt/radar/api/agent-review.js
NoNewPrivileges=true
PrivateTmp=true
