# motor/radar-varredura.service
[Unit]
Description=Radar Fundiario - varredura noturna de anuncios (fonte A)
After=network.target postgresql.service ollama.service

[Service]
Type=oneshot
WorkingDirectory=/opt/radar/api
EnvironmentFile=/opt/radar/api/.env
ExecStart=/usr/bin/node /opt/radar/api/varredura.js
User=www-data
NoNewPrivileges=true
PrivateTmp=true
TimeoutStartSec=4h
